# launchd/dev.rocky.embedding.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.rocky.embedding</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/steve/amaze_s3/rocky/.venv/bin/rocky</string>
    <string>embed</string>
    <string>qwen3-embed-4b</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/steve/amaze_s3/rocky/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/steve</string>
    <key>ROCKY_HOST</key>
    <string>127.0.0.1</string>
    <key>ROCKY_EMBED_PORT</key>
    <string>7778</string>
    <key>ROCKY_RUNTIME_ROOT</key>
    <string>/Users/steve/amaze_s3/rocky/.rocky</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/steve/amaze_s3/rocky</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/steve/amaze_s3/rocky/.rocky/logs/rocky-embedding.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/steve/amaze_s3/rocky/.rocky/logs/rocky-embedding.err.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
